# kojicron.timer
[Unit]
Description=Daily koji repo regeneration
Wants=network-online.target

[Timer]
OnCalendar = *-*-* 04:00:00

[Install]
WantedBy = timers.target
